# batterymonitor-service.conf
description "BatteryMonitor"
author "Thomas Crha"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
  
respawn
exec python /home/BatteryMonitor/BatteryMonitor.py start
